# Docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisor.log

[program:configservice]
command=/home/microservice-baseline/start-service.sh
user=microservice-baseline
directory=/home/microservice-baseline
autorestart = true

[program:cron]
command         = /usr/sbin/cron -f
stdout_logfile	= /var/log/supervisor/%(program_name)s.log
stderr_logfile	= /var/log/supervisor/%(program_name)s.log
autorestart	    = true